# tools/bug-scanner/sml-bug-scanner.service
[Unit]
Description=SML bug-scanner — 掃描 bug 回報區、交 claude 判定並貼回報告
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=smlbot
WorkingDirectory=/opt/sml/repo
Environment=HOME=/home/smlbot
Environment=AWS_REGION=ap-southeast-1
# kill switch:改 0 即停用(記得 daemon-reload 不需要,env 每次啟動讀取)
Environment=BUG_SCANNER_ENABLED=1
# 0=只診斷+附 patch(預設,安全);1=明確 bug 就在 bug-scan/<日期> 分支 commit。
# ⚠️ 暫關:多數 SML bug 在別的 repo(sweetbot-next/score-repo),本 worktree 動不到,
# 開了 claude 反而會在本 repo 亂刪無關檔。需更嚴的圍欄才重開。
Environment=BUG_AUTO_FIX=0
# 診斷報告貼到「維修頻道」bug回報專區(開發中心,客戶看不到);入口頻道 BUG_CHANNEL_ID 保持乾淨給客戶回報。
Environment=BUG_POST_CHANNEL=1527655549532115054
ExecStart=/usr/bin/python3 /opt/sml/repo/tools/bug-scanner/scan.py
TimeoutStartSec=1800
Nice=10

[Install]
WantedBy=multi-user.target
